# flatpak-updater.timer
[Unit]
Description=Run Flatpak Updates Daily
Requires=flatpak-updater.service

[Timer]
OnCalendar=*-*-* 22:30:00

[Install]
WantedBy=timers.target
